# ansible/templates/quadlet/home-assistant.container.j2
[Unit]
Description=Home Assistant
Wants=network-online.target
After=network-online.target

[Container]
Image=ghcr.io/home-assistant/home-assistant:stable
ContainerName=home-assistant
# Host networking is the upstream-recommended mode: HA needs it for mDNS,
# SSDP, and DHCP discovery of integrations (Sonos, Cast, Plex, etc.). The
# UI listens on :8123 directly on the host.
Network=host
# Run HA as host josh (UID 1000) instead of container-root, so DBus
# EXTERNAL auth (which compares the container's getuid() against the host
# SO_PEERCRED on the bus socket) succeeds. Without this combo, the BT
# integration hits 'REJECTED: [EXTERNAL]' against BlueZ. /config is
# josh-owned on the host so file ownership lines up.
UserNS=keep-id
User=1000:1000
# Drop SELinux container_t confinement so HA can connect to the host's
# system DBus socket (for BlueZ). Bazzite's default container_t silently
# (dontaudit) denies AF_UNIX connect() to system_dbusd_var_run_t.
SecurityLabelDisable=true
# NET_ADMIN + NET_RAW: needed for HA's DHCP integration discovery
# (aiodhcpwatcher opens an AF_PACKET socket) and BT adapter recovery.
AddCapability=NET_ADMIN
AddCapability=NET_RAW
Environment=TZ=America/New_York
Volume={{ infra_dir }}/home-assistant:/config
Volume=/etc/localtime:/etc/localtime:ro
# Host system DBus, for HA -> BlueZ. Note this still hits the rootless
# DBus EXTERNAL auth mismatch (container UID 0 vs host SO_PEERCRED 1000)
# unless paired with a non-root HA user + UserNS=keep-id.
Volume=/run/dbus:/run/dbus:ro
# Zigbee coordinator is a network-attached SMLIGHT SLZB (Z-Stack/znp),
# reached over TCP at socket://<slzb-ip>:6638 — configured in ZHA, not
# here. No USB device passthrough is needed. (Previously this passed
# AddDevice=/dev/zigbee for a USB Sonoff dongle; that dongle was removed
# when the coordinator moved onto the LAN. Passing a non-existent device
# makes `podman run` fail with exit 125, so it must stay absent unless a
# USB radio is reintroduced.)

[Service]
Restart=on-failure
TimeoutStartSec=900

[Install]
WantedBy=default.target
