# ansible/templates/quadlet/pihole.container.j2
[Unit]
Description=Pi-hole DNS + ad-blocker
Wants=network-online.target
After=network-online.target

[Container]
Image=docker.io/pihole/pihole:latest
ContainerName=pihole
Network=arr.network
# Bind to bazz's LAN IP only so we don't collide with systemd-resolved's
# 127.0.0.53:53 stub. Make {{ pihole_host_ip }} stable via a DHCP
# reservation on the router (that's the "fixed IP" clients point at).
PublishPort={{ pihole_host_ip }}:53:53/tcp
PublishPort={{ pihole_host_ip }}:53:53/udp
PublishPort={{ pihole_host_ip }}:8053:80/tcp
Environment=TZ=America/New_York
Environment=FTLCONF_dns_listeningMode=all
EnvironmentFile={{ infra_dir }}/pihole/pihole.env
Volume={{ infra_dir }}/pihole/etc-pihole:/etc/pihole:Z
Volume={{ infra_dir }}/pihole/etc-dnsmasq.d:/etc/dnsmasq.d:Z

[Service]
Restart=on-failure
TimeoutStartSec=900

[Install]
WantedBy=default.target
